# services/qbittorrent/run
#!/bin/bash

exec /sbin/setuser abc qbittorrent --webui-port=$PORT
