# services/qbittorrent-nox/run
#!/bin/bash

exec /sbin/setuser abc qbittorrent-nox --webui-port=$PORT
